# share/nti/dataserver/etc/supervisord.conf

[supervisord]
childlogdir = /home/ntiuser/buildout/var/log
logfile = /home/ntiuser/buildout/var/log/supervisord.log
logfile_maxbytes = 0
logfile_backups = 10
loglevel = info
pidfile = /home/ntiuser/buildout/var/supervisord.pid
umask = 022
nodaemon = false
nocleanup = True


environment = DATASERVER_DIR=/home/ntiuser/buildout,PYTHONHASHSEED=random

[inet_http_server]
port = 127.0.0.1:9001
username = 
password = 

[supervisorctl]
serverurl = http://127.0.0.1:9001
username = 
password = 

[ctlplugin:laforge]
supervisor.ctl_factory = mr.laforge.controllerplugin:make_laforge_controllerplugin

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[rpcinterface:laforge]
supervisor.rpcinterface_factory=mr.laforge.rpcinterface:make_laforge_rpcinterface

[program:pserve]
command = /home/ntiuser/buildout/bin/nti_pserve /home/ntiuser/buildout/etc/pserve.ini
process_name = pserve
directory = /home/ntiuser/buildout/bin
priority = 999
redirect_stderr = false
stopwaitsecs = 30
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:solr_indexer]
command = /home/ntiuser/buildout/bin/nti_solr_indexer -v --slugs
process_name = solr_indexer
directory = /home/ntiuser/buildout/bin
priority = 99
redirect_stderr = false
autostart = false
stopwaitsecs = 30
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:metadata]
command = /home/ntiuser/buildout/bin/nti_metadata_processor -v --slugs
process_name = metadata
directory = /home/ntiuser/buildout/bin
priority = 99
redirect_stderr = false
autostart = false
stopwaitsecs = 30
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:library_renderer]
command = /home/ntiuser/buildout/bin/nti_library_renderer -v -r 1 -s 2 --slugs
process_name = library_renderer
directory = /home/ntiuser/buildout/bin
priority = 99
redirect_stderr = false
autostart = false
stopwaitsecs = 30
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:scheduled_job_dispatcher]
command = /home/ntiuser/buildout/bin/nti_scheduled_dispatcher -v --slugs
process_name = scheduled_job_dispatcher
directory = /home/ntiuser/buildout/bin
priority = 99
redirect_stderr = false
autostart = false
stopwaitsecs = 30
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:scheduled_job_executor]
command = /home/ntiuser/buildout/bin/nti_scheduled_executor -v --slugs
process_name = scheduled_job_executor
directory = /home/ntiuser/buildout/bin
priority = 99
redirect_stderr = false
autostart = false
stopwaitsecs = 30
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[program:qp_mailer]
command = /home/ntiuser/buildout/bin/nti_mailer_qp_process -v /home/ntiuser/buildout/var/maildir
process_name = qp_mailer
directory = /home/ntiuser/buildout/bin
priority = 99
redirect_stderr = false
stopwaitsecs = 30
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0

[eventlistener:crashmail_pserve]
command = /home/ntiuser/buildout/bin/crashmail -p pserve -m server-errors@nextthought.com -o Error 
events = PROCESS_STATE_EXITED
process_name=crashmail_pserve
environment=SUPERVISOR_USERNAME='',SUPERVISOR_PASSWORD='',SUPERVISOR_SERVER_URL='http://127.0.0.1:9001'


[eventlistener:crashmail_scheduled_job_executor]
command = /home/ntiuser/buildout/bin/crashmail -p scheduled_job_executor -m server-errors@nextthought.com -o Error 
events = PROCESS_STATE_EXITED
process_name=crashmail_scheduled_job_executor
environment=SUPERVISOR_USERNAME='',SUPERVISOR_PASSWORD='',SUPERVISOR_SERVER_URL='http://127.0.0.1:9001'


[eventlistener:crashmail_scheduled_job_dispatcher]
command = /home/ntiuser/buildout/bin/crashmail -p scheduled_job_dispatcher -m server-errors@nextthought.com -o Error 
events = PROCESS_STATE_EXITED
process_name=crashmail_scheduled_job_dispatcher
environment=SUPERVISOR_USERNAME='',SUPERVISOR_PASSWORD='',SUPERVISOR_SERVER_URL='http://127.0.0.1:9001'


[eventlistener:fatalmail]
command = /home/ntiuser/buildout/bin/fatalmailbatch -t server-errors@nextthought.com -f automated.emails@nextthought.com -s Error -H :587 -u  -p 
events = PROCESS_STATE,TICK_60
process_name=fatalmail
environment=SUPERVISOR_USERNAME='',SUPERVISOR_PASSWORD='',SUPERVISOR_SERVER_URL='http://127.0.0.1:9001'


[group:server]
programs = pserve,metadata,solr_indexer,library_renderer,scheduled_job_dispatcher,scheduled_job_executor
priority = 999
